# launchd/net.medarov.docspell.verify.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent: runs verify_docspell.py --brief once a day at 09:00.

  Install:
      cp launchd/net.medarov.docspell.verify.plist ~/Library/LaunchAgents/
      launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/net.medarov.docspell.verify.plist
      # to disable:
      launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/net.medarov.docspell.verify.plist

  Logs:
      ~/Library/Logs/docspell-verify.log

  Note: the agent reads the Docspell password from macOS keychain via
  the DOCSPELL_PASSWORD env var, which the wrapper script populates.
  Make sure you have run:
      security add-generic-password -s docspell -a library/dmedarov -w
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>net.medarov.docspell.verify</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>DOCSPELL_PASSWORD=$(security find-generic-password -s docspell -a library/dmedarov -w) /opt/homebrew/bin/python3 /Users/dmedarov/CODING/Docspell/verify_docspell.py --brief</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/dmedarov/CODING/Docspell</string>

  <key>StandardOutPath</key>
  <string>/Users/dmedarov/Library/Logs/docspell-verify.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/dmedarov/Library/Logs/docspell-verify.log</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
